# overlay/usr/local/etc/rc.d/deluged_restart
#!/bin/sh
#
# PROVIDE: deluged_restart
# REQUIRE: deluged openvpn ipfw
#

. /etc/rc.subr

name="deluged_restart"
start_cmd="${name}_start"
stop_cmd=":"
